# smartproxy.service.in
[Unit]
Description=SmartProxy
DefaultDependencies=no

[Service]
ExecStart=@prefix@/bin/smartproxy

[Install]
WantedBy=multi-user.target
