# smartproxy.service.in
[Unit]
Description=SmartProxy
DefaultDependencies=no

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/smartproxy
LimitNOFILE=900000

[Install]
WantedBy=multi-user.target
