# smartproxy.service.in
[Unit]
Description=SmartProxy
DefaultDependencies=no

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/smartproxy

[Install]
WantedBy=multi-user.target
